# build/md-feed-01/systemd-feed.service
[Unit]
Description=Synthetic Market Data Feed

[Service]
ExecStart=/usr/bin/python3 feed/generator.py
Restart=always

[Install]
WantedBy=multi-user.target